# .config/systemd/user/vdirsyncer-sync.service
[Unit]
Description=Sync CalDAV/CardDAV via vdirsyncer

[Service]
Type=oneshot
ExecStart=/usr/bin/vdirsyncer sync

